# cron.hourly/tzupdate
#!/bin/dash

# Run tzupdate
tzupdate 2>&1 1>/dev/null
