# profile.d/powerlevel10k.sh
export POWERLEVEL10K_INSTALLATION_PATH="/usr/share/zsh-theme-powerlevel10k"
