# systemd/system/sbp-tzupdate.timer
[Unit]
Description=Run tzupdate every hour

[Timer]
OnBootSec=5min

[Install]
WantedBy=timers.target
