# systemd/system/hostapd.service.d/override.conf
[Unit]
BindsTo=sys-subsystem-net-devices-wifi\x2dwork.device
After=sys-subsystem-net-devices-wifi\x2dwork.device