# systemd/system/getty@tty1.service.d/sbp.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin sbp --noclear %I $TERM
